# src/disks/45h/etc/sysinit/system/poweroff.target
[Unit]
Description=Shutdown of system
Conflicts=reboot.target halt.target

[Service]